# systemd/ps-factory-daily-radar.timer
[Unit]
Description=Run PS Factory Daily Radar at 07:30 UTC

[Timer]
OnCalendar=*-*-* 07:30:00 UTC
Persistent=true
Unit=ps-factory-daily-radar.service

[Install]
WantedBy=timers.target
